# service.cron
* * * * * ./run_service 5 >> ~/minimoto.log
